# resources/supervisor/ss-local.conf
[program:ss-local]
directory=/opt/ss-local
command=/opt/ss-local/scripts/start-ss-local.sh
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/ss-local.stdout.log
stderr_logfile=/var/log/ss-local.stderr.log
environment=PHP_BIN="/usr/bin/php",SS_SERVER="your-node.example.com",SS_PORT="18001",SS_CIPHER="aes-256-gcm",SS_PASSWORD="your-password",SS_UDP="0",SS_LISTEN="127.0.0.1:1080",SS_WORKER_COUNT="1",SS_MAX_CONNECTIONS="1024",SS_ALLOW_IPS="127.0.0.1,::1",SS_CONNECT_TIMEOUT="10",SS_CONNECT_RETRIES="1",SS_RETRY_DELAY_MS="250",SS_IDLE_TIMEOUT="900",SS_MAX_SEND_BUFFER="4194304",SS_STATUS_FILE="/var/run/ss-local.status.json",SS_STATUS_INTERVAL="10"
